# templates/systemd-adcollection.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Run cartography AD collection daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=adcollection.target
