# templates/systemd-cartography@.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Run cartography daily - %i

[Timer]
OnCalendar={{ cartography_timer_oncalendar }}
Persistent=true

[Install]
WantedBy=cartography@%i.target
