# templates/systemd-cartography.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Run cartography daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=cartography.target
